# data/51-automedia.rules
ACTION=="change", SUBSYSTEM=="block", RUN+="/opt/udfbackup/bash_wrapper.sh %k"