# data/51-automedia.rules
ACTION=="change", SUBSYSTEM=="block", RUN+="/bin/su -l -c '/opt/udfbackup/bash_wrapper.sh %k' -s /bin/bash pragma"